# test/integration/suites/ghostunnel-federation/conf/downstream/supervisord/supervisord.conf
[supervisord]
nodaemon=true
loglevel=debug

[unix_http_server]
file = /tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl = unix:///tmp/supervisor.sock

[program:spire-agent]
command = /opt/spire/bin/spire-agent run -config /opt/spire/conf/agent/agent.conf

[program:ghostunnel]
command = /usr/bin/ghostunnel @/opt/ghostunnel/ghostunnel.flags

[program:socat]
command = /usr/bin/socat -d -d TCP-LISTEN:8000,fork TCP:localhost:8001
