# test/integration/cassandra-suites/envoy-sds-v3/conf/supervisord.conf
[supervisord]
nodaemon=true
loglevel=debug

[program:spire-agent]
command = /opt/spire/bin/spire-agent run -config /opt/spire/conf/agent/agent.conf

[program:envoy]
command = /usr/local/bin/envoy -l debug -c /opt/envoy/conf/envoy.yaml
